# tools/sovereign-chatgpt-mcp/deploy/sovereign-openai-tunnel.service
[Unit]
Description=Sovereign OpenAI Secure MCP Tunnel
After=network-online.target docker.service sovereign-chatgpt-broker.service
Wants=network-online.target
Requires=sovereign-chatgpt-broker.service

[Service]
Type=simple
User=sovereign-tunnel
Group=sovereign-tunnel
EnvironmentFile=/opt/sovereign-chatgpt-tools/tunnel.env
Environment=HOME=/var/lib/sovereign-tunnel
WorkingDirectory=/var/lib/sovereign-tunnel
# Use the shared protocol-aware health checker. It builds the JSON-RPC payload
# in Python, validates the full initialize response and avoids systemd/shell
# quoting corruption that previously caused intermittent HTTP 400 responses.
ExecStartPre=/usr/bin/test -x /usr/local/bin/tunnel-client
ExecStartPre=/usr/bin/python3 /opt/sovereign-chatgpt-tools/mcp_protocol_health.py --url http://127.0.0.1:8090/mcp --wait-seconds 20 --timeout-seconds 5
ExecStart=/usr/local/bin/tunnel-client run --profile ${TUNNEL_PROFILE}
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
StateDirectory=sovereign-tunnel
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/opt/sovereign-chatgpt-tools/tunnel.env
ReadWritePaths=/var/lib/sovereign-tunnel
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target

